# systemd_drafts/tokenoskobi-phase9-observation-runtime.timer
# Draft only. Not installed/enabled by PHASE9C.
[Unit]
Description=Tokenoskobi Phase9 Commercial Observation Runtime Timer Draft

[Timer]
OnBootSec=15min
OnUnitActiveSec=15min
AccuracySec=2min
Persistent=false

[Install]
WantedBy=timers.target
